# video_protection_bot.service
[Unit]
Description=Telegram Video Protection Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/azprotectionbot
ExecStart=/opt/azprotectionbot/.venv/bin/python /opt/azprotectionbot/bot.py
Restart=always
RestartSec=5
EnvironmentFile=/opt/azprotectionbot/.env

# Hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
